# script/vimtricks.service
[Unit]
Description=vimtricks

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=vimtricks

[Install]
WantedBy=multi-user.target
